# setup/etc/cron.d/renew_pass_token
SHELL=/bin/bash
MAILTO=""

55 4 * * *	root	/bin/bash /root/bin/init_pass.sh
